# php-workers/conf.d/worker.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=100

[program:yii-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/var/www/app.local/yii queue/listen --verbose=1 --color=0
autostart=true
autorestart=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/app.local/runtime/logs/yii-queue-worker.log
priority=200
